# root/etc/init.d/clash-meta
#!/bin/sh /etc/rc.common
# Clash.Meta init script (placeholder)